# modules/puppetmaster/files/update-fastly-ips.cron
MAILTO=root
@daily	root	update-fastly-ips /srv/puppet.debian.org/puppet-facts/fastly_ranges.yaml
